# rootfs/etc/s6-overlay/s6-rc.d/init-authelia/run
#!/command/with-contenv bashio
# ==============================================================================
# Benoit Anastay Add-on: Authelia
# Runs the Authelia daemon
# ==============================================================================

bashio::log.info "Starting Authelia...."

exec authelia --config /share/authelia/config/configuration.yml
