# authelia/rootfs/etc/s6-overlay/s6-rc.d/authelia/run
#!/command/with-contenv bashio
# ==============================================================================
# Benoit Anastay Add-on: Authelia
# Runs the Authelia daemon
# ==============================================================================

bashio::log.info "Starting Authelia...."

# Support for https://www.authelia.com/configuration/methods/files/#file-filters.
# For compatibility/stability reasons use the variable instead of the command-line argument:
# "...the name of the CLI argument will change (we suggest using the environment variable which will not)"
export X_AUTHELIA_CONFIG_FILTERS="$(bashio::config 'config_filters' '')"

exec authelia --config /config/configuration.yml
